# perfsonar-archive/perfsonar-archive/config/systemd/opensearch.service.d/override.conf
[Unit]
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Restart=always
RestartSec=10